# deploy/wirenest.service
[Unit]
# Reference template. The one-command installer (deploy/install.sh) generates
# its own copy with a random admin password and the detected public IP filled in.
Description=WireNest - WireGuard management panel
StartLimitIntervalSec=0
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as root: the panel drives wg-quick (interface up/down/restart), systemctl
# (boot autostart toggle) and sysctl (IP forwarding) — all privileged operations.
Environment=WIRENEST_ADDR=:8000
Environment=WIRENEST_DATA_DIR=/var/lib/wirenest
Environment=WIRENEST_WG_CONF=/etc/wireguard/wg0.conf
# Public host clients dial (optional; settable in the panel too).
Environment=WIRENEST_ENDPOINT=
# CHANGE THESE before exposing the panel.
Environment=WIRENEST_ADMIN_USER=admin
Environment=WIRENEST_ADMIN_PASS=admin

ExecStart=/usr/local/bin/wirenest-panel
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
